# recipes-common/recipes-extend/v4l2-init/files/v4l2-init.service
[Unit]
Description=VIN initializing setting

[Service]
ExecStart=/home/root/vin-init.sh

[Install]
WantedBy=multi-user.target
